# debian/09console-setup/10-console.rules
# Console initialization - keyboard, font, etc.
KERNEL=="tty0",		RUN+="/sbin/initqueue --onetime --unique --name console-setup /lib/udev/console-setup"
